# deploy/systemd/user/memoryos-api.service
[Unit]
Description=MemoryOS API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/VaultMind
EnvironmentFile=%h/.config/memoryos.env
ExecStart=/usr/bin/python3 -m memoryos.app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
